# packaging/dot-ddns.install
post_install() {
  cat <<'EOF'
Create config with:
  sudo dot-ddns init --domain one.one.one.one --bootstrap 9.9.9.9,1.1.1.1

Then enable the service with:
  sudo systemctl enable --now dot-ddns.service
EOF
}

post_upgrade() {
  post_install
}
